# discord/opteee-bot.service
[Unit]
Description=OPTEEE Discord Bot
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/gitrepos/opteee/discord
Environment=PYTHONPATH=%h/gitrepos/opteee/discord
ExecStart=%h/gitrepos/opteee/venv/bin/python discord_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target